# launchd/ai.hermes.discord-voice-transcriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.hermes.discord-voice-transcriber</string>
  <key>Comment</key>
  <string>Discord VC recorder/transcriber for nikechan-discord-public</string>
  <key>UserName</key>
  <string>nikenike</string>
  <key>GroupName</key>
  <string>staff</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>node</string>
    <string>/Users/nikenike/.hermes/profiles/nikechan-discord-public/bin/discord-voice-transcriber.mjs</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/nikenike/.hermes/profiles/nikechan-discord-public</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/nikenike</string>
    <key>HERMES_HOME</key>
    <string>/Users/nikenike/.hermes/profiles/nikechan-discord-public</string>
    <key>PATH</key>
    <string>/Users/nikenike/.hermes/hermes-agent/node_modules/.bin:/usr/local/node/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>DISCORD_VOICE_CHANNEL_ID</key>
    <string>1452811457925480580</string>
    <key>DISCORD_VOICE_START_GRACE_SECONDS</key>
    <string>30</string>
    <key>DISCORD_VOICE_END_GRACE_SECONDS</key>
    <string>30</string>
    <key>DISCORD_VOICE_TRANSCRIBE_PROVIDER</key>
    <string>local</string>
    <key>DISCORD_VOICE_LOCAL_WHISPER_MODEL</key>
    <string>models/whisper/ggml-large-v3-turbo.bin</string>
    <key>DISCORD_VOICE_LOCAL_WHISPER_THREADS</key>
    <string>4</string>
    <key>DISCORD_VOICE_LOCAL_MAX_SESSION_SECONDS</key>
    <string>10800</string>
    <key>DISCORD_VOICE_LOCAL_MAX_TOTAL_SESSION_SECONDS</key>
    <string>43200</string>
    <key>DISCORD_VOICE_LOCAL_SEGMENT_SECONDS</key>
    <string>1200</string>
    <key>DISCORD_VOICE_LOCAL_MAX_SESSION_CHUNKS</key>
    <string>10000</string>
    <key>DISCORD_VOICE_MAX_SESSION_SECONDS</key>
    <string>1800</string>
    <key>DISCORD_VOICE_MAX_SESSION_CHUNKS</key>
    <string>250</string>
    <key>DISCORD_VOICE_MAX_SESSION_AUDIO_BYTES</key>
    <string>20971520</string>
    <key>DISCORD_VOICE_MIN_SESSION_AUDIO_BYTES</key>
    <string>8192</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/nikenike/.hermes/profiles/nikechan-discord-public/logs/voice-transcriber.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nikenike/.hermes/profiles/nikechan-discord-public/logs/voice-transcriber.error.log</string>
</dict>
</plist>
